# apps/review-prototype/deploy/chatcommons-review.service
[Unit]
Description=ChatCommons isolated UI review API
After=network.target

[Service]
Type=simple
User=www
Group=www
UMask=0027
WorkingDirectory=/var/www/keyanruanjian.com/chatcommons-review/current
EnvironmentFile=/etc/chatcommons-review.env
ExecStart=/usr/bin/python3 /var/www/keyanruanjian.com/chatcommons-review/current/server.py
Restart=on-failure
RestartSec=3
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
NoNewPrivileges=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/var/lib/chatcommons-review

[Install]
WantedBy=multi-user.target
